# Components/motd/update-motd.d/10-uname
#!/bin/bash
description="$(lsb_release -ds)"
release="$(lsb_release -rs)"
echo "Welcome to ${description% *} ${release^} ("$(uname -r) $(uname -m)")"
