# Components/motd/update-motd.d/10-uname
#!/bin/sh
echo "Welcome to Debian GNU/Linux Unstable" "("$(uname -r) $(uname -m)")"
