# Components/motd/update-motd.d/11-custom-motd
#!/bin/bash
printf "[1;32m                     __                         \033[1;31m__ \n"
printf "[1;32m   _________  ____ _/ /_____  \033[1;31m_________  ____  / /_\n"
printf "[1;32m  / ___/ __ \/ __ \`/ //_/ _ \[1;31m/ ___/ __ \/ __ \/ __/\n"
printf "[1;32m (__  ) / / / /_/ / ,< /  __/ \033[1;31m/  / /_/ / /_/ / /_\n"
printf "[1;32m/____/_/ /_/\__,_/_/|_|\___/\033[1;31m_/   \____/\____/\__/\n"
                                                   
printf "[0m"
