# Components/motd/update-motd.d/98-reboot-required
#!/bin/sh

[ -f /var/run/reboot-required ] && cat /var/run/reboot-required
